# etc/crontab
MAILTO=""
*/5 * * * * /usr/bin/php <path to ttrss reader update script> --feeds --quiet
